# src/ouroborOS-profile/airootfs/etc/systemd/system/bluetooth.service.d/experimental.conf
# bluetooth.service.d/experimental.conf
#
# Drop-in for bluetooth.service that enables experimental BlueZ mode.
#
# WHY --experimental:
#   BlueZ experimental mode activates D-Bus APIs that are still being
#   stabilised upstream but are fully functional and safe:
#
#   - AdvertisingMonitor API (org.bluez.AdvertisementMonitor1)
#       Used by Chrome/Chromium (≥ 108) and Firefox (≥ 120) on Linux to scan
#       for BLE advertisements from phone authenticators during CTAP2 hybrid
#       transport (the "QR code passkey" flow). Without this API the browser
#       falls back to showing "Bluetooth is not available on this computer."
#
#   - FIDO2 BLE GATT enumeration improvements
#       More reliable discovery of BLE FIDO2 tokens exposing the FIDO2
#       GATT service (UUID 0x0000FFFD-0000-1000-8000-00805F9B34FB).
#
#   - BlueZ mesh, BLE audio (if applicable)
#
# SECURITY NOTE:
#   The experimental APIs expose no additional attack surface beyond
#   what bluetooth.service already exposes. Only local D-Bus clients
#   (typically the browser or our-fido2) can use these APIs.
#
# See also: /etc/bluetooth/main.conf (Experimental = true for D-Bus API level)

[Service]
# Clear the default ExecStart to override it with --experimental.
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd --experimental
