# src/ouroborOS-profile-aarch64/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
# Auto-login as root on tty1 in the live environment
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
