# src/ouroborOS-profile/airootfs/etc/systemd/system/our-box-autostart.service
# our-box-autostart.service — auto-start systemd-nspawn containers at boot
#
# Reads /etc/our-box/autostart.conf and runs `our-box start <name>` for each
# container listed (one per line). Empty lines and # comments are ignored.
#
# The service runs as a oneshot unit after systemd-machined is available.
# Individual container start failures are logged but do NOT abort remaining
# entries — we want best-effort behaviour so one broken container doesn't
# block the rest.
#
# Enable:  systemctl enable our-box-autostart.service
# Disable: systemctl disable our-box-autostart.service

[Unit]
Description=Auto-start our-box containers at boot
After=systemd-machined.service
Wants=systemd-machined.service
ConditionPathExists=/etc/our-box/autostart.conf

[Service]
Type=oneshot
ExecStart=/usr/local/bin/our-box-autostart

[Install]
WantedBy=multi-user.target
