# src/ouroborOS-profile/airootfs/etc/systemd/system/our-snapshot-prune.timer
# our-snapshot-prune.timer — daily trigger for snapshot pruning
#
# Enable: systemctl enable our-snapshot-prune.timer
# The timer fires once daily at 03:30 with a 15-minute randomized delay
# to avoid thundering-herd on systems that boot at the same time.

[Unit]
Description=Daily ouroborOS snapshot pruning

[Timer]
OnCalendar=daily
RandomizedDelaySec=15min
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
