# src/ouroborOS-profile-aarch64/airootfs/etc/systemd/system/ouroboros-homed-migration.service
[Unit]
Description=ouroborOS systemd-homed migration
Documentation=https://github.com/Arkhur-Vo/ouroborOS
DefaultDependencies=no
After=systemd-homed.service
After=local-fs.target
Before=getty.target
Before=display-manager.service
Wants=systemd-homed.service
ConditionPathExists=/etc/ouroboros/homed-migration-pending
ConditionPathExists=/etc/ouroboros/homed-migration.conf

[Service]
Type=oneshot
EnvironmentFile=/etc/ouroboros/homed-migration.conf
ExecStart=/usr/local/lib/ouroboros/homed-migrate.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
